# roles/k3s/templates/k3s-cni-bridge-fixer.service.j2
# Managed by Ansible — do not edit manually.
# Workaround for k3s embedded CNI skipping cni0 bridge attachment.
[Unit]
Description=k3s CNI bridge fixer — attaches pod veths to cni0
After=k3s.service
BindsTo=k3s.service

[Service]
Type=simple
ExecStart=/usr/local/bin/k3s-cni-bridge-fixer.sh
Restart=always
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=k3s-cni-fixer

[Install]
WantedBy=multi-user.target
