# roles/ansible-pull/templates/ansible-pull.timer.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=Ansible Pull Timer — run every {{ ansible_pull_interval }}
Documentation=https://github.com/nwlucas/olympus-sdk

[Timer]
OnBootSec=5min
OnUnitActiveSec={{ ansible_pull_interval }}
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
